# systemd/service-timer-files/hello-world-socket.service
[Unit]
Description=Hello-World-Socket

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target