# systemd/timer/remove-chrome-policy.timer
[Unit]
Description=Runs command every day 

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target